# packaging/deb/prerm
#!/usr/bin/env bash
set -euo pipefail

LOG_FILE="/var/log/wps-read-aloud-install.log"
mkdir -p /var/log
exec > >(tee -a "${LOG_FILE}") 2>&1
echo "[$(date '+%Y-%m-%d %H:%M:%S')] prerm $*"

SERVICE_NAME="@SERVICE_NAME@"

if [[ "${1:-}" = "remove" || "${1:-}" = "deconfigure" ]]; then
  systemctl disable --now "${SERVICE_NAME}" || true
fi
